# build/docker/gorevoke.service
[Unit]
Description=GoRevoke systemd service unit file.

[Service]
ExecStart=/usr/local/bin/gorevoke/gorevoke

[Install]
WantedBy=multi-user.target